# config/crontab
0 */4 * * * /usr/local/bin/python /app/manage.py fetch_events
